# 101-ps_events_windows.xml
<group name="windows,powershell,ad_recon,">

  <!-- ============================================================
       SECTION 1 — Basic System Reconnaissance
       These are the first commands attackers run after initial access
       ============================================================ -->

  <!-- whoami variants - identity discovery -->
  <rule id="100400" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)whoami</field>
    <options>no_full_log</options>
    <description>PS Recon: Identity discovery via whoami on $(win.system.computer) by $(win.eventdata.path)</description>
    <mitre>
      <id>T1033</id>
    </mitre>
    <group>recon,discovery,</group>
  </rule>

  <!-- System info gathering -->
  <rule id="100401" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)systeminfo|Get-ComputerInfo|Get-WmiObject Win32_ComputerSystem|Get-CimInstance Win32_ComputerSystem</field>
    <options>no_full_log</options>
    <description>PS Recon: System information gathering on $(win.system.computer)</description>
    <mitre>
      <id>T1082</id>
    </mitre>
    <group>recon,discovery,</group>
  </rule>

  <!-- Network configuration discovery -->
  <rule id="100402" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)ipconfig|Get-NetIPAddress|Get-NetIPConfiguration|Get-NetAdapter</field>
    <options>no_full_log</options>
    <description>PS Recon: Network configuration discovery on $(win.system.computer)</description>
    <mitre>
      <id>T1016</id>
    </mitre>
    <group>recon,discovery,network_recon,</group>
  </rule>

  <!-- Running process enumeration -->
  <rule id="100403" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-Process|tasklist|Get-WmiObject Win32_Process|Get-CimInstance Win32_Process</field>
    <options>no_full_log</options>
    <description>PS Recon: Running process enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1057</id>
    </mitre>
    <group>recon,discovery,</group>
  </rule>

  <!-- Installed software enumeration -->
  <rule id="100404" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-WmiObject Win32_Product|Get-ItemProperty HKLM.*Uninstall|winrm enumerate</field>
    <options>no_full_log</options>
    <description>PS Recon: Installed software enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1518</id>
    </mitre>
    <group>recon,discovery,</group>
  </rule>

  <!-- Local user and group enumeration -->
  <rule id="100405" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-LocalUser|Get-LocalGroup|Get-LocalGroupMember|net localgroup|net user</field>
    <options>no_full_log</options>
    <description>PS Recon: Local user/group enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1087.001</id>
    </mitre>
    <group>recon,discovery,account_discovery,</group>
  </rule>

  <!-- Open ports and network connections -->
  <rule id="100406" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-NetTCPConnection|Get-NetUDPEndpoint|netstat</field>
    <options>no_full_log</options>
    <description>PS Recon: Network connection/port enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1049</id>
    </mitre>
    <group>recon,discovery,network_recon,</group>
  </rule>

  <!-- Scheduled tasks enumeration -->
  <rule id="100407" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-ScheduledTask|schtasks</field>
    <options>no_full_log</options>
    <description>PS Recon: Scheduled task enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>recon,discovery,</group>
  </rule>

  <!-- Services enumeration -->
  <rule id="100408" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-Service|sc query|Get-WmiObject Win32_Service|Get-CimInstance Win32_Service</field>
    <options>no_full_log</options>
    <description>PS Recon: Windows service enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1007</id>
    </mitre>
    <group>recon,discovery,</group>
  </rule>

  <!-- ============================================================
       SECTION 2 — Active Directory Enumeration
       High value targets - attackers map AD before moving laterally
       ============================================================ -->

  <!-- AD user enumeration -->
  <rule id="100410" level="10">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-ADUser|Get-ADObject|([adsisearcher])</field>
    <options>no_full_log</options>
    <description>PS Recon: Active Directory user enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1087.002</id>
    </mitre>
    <group>recon,discovery,ad_enumeration,</group>
  </rule>

  <!-- AD group enumeration -->
  <rule id="100411" level="10">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-ADGroup|Get-ADGroupMember|Get-ADPrincipalGroupMembership</field>
    <options>no_full_log</options>
    <description>PS Recon: Active Directory group enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1069.002</id>
    </mitre>
    <group>recon,discovery,ad_enumeration,</group>
  </rule>

  <!-- AD computer enumeration -->
  <rule id="100412" level="10">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-ADComputer|Get-ADDomainController|Get-ADDomain|Get-ADForest</field>
    <options>no_full_log</options>
    <description>PS Recon: Active Directory domain/computer enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1018</id>
    </mitre>
    <group>recon,discovery,ad_enumeration,</group>
  </rule>

  <!-- AD trust enumeration -->
  <rule id="100413" level="10">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-ADTrust|nltest /domain_trusts|Get-NetDomainTrust</field>
    <options>no_full_log</options>
    <description>PS Recon: Active Directory trust enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1482</id>
    </mitre>
    <group>recon,discovery,ad_enumeration,</group>
  </rule>

  <!-- SPN / Kerberoastable account enumeration -->
  <rule id="100414" level="12">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-ADUser.*ServicePrincipalName|setspn|KerberosRequestorSecurityToken|Request-SPNTicket</field>
    <options>no_full_log</options>
    <description>PS Recon: SPN / Kerberoastable account enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1558.003</id>
    </mitre>
    <group>recon,discovery,ad_enumeration,kerberoasting,</group>
  </rule>

  <!-- GPO and OU enumeration -->
  <rule id="100415" level="10">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-GPO|Get-ADOU|Get-ADOrganizationalUnit|Get-GPPermission|gpresult</field>
    <options>no_full_log</options>
    <description>PS Recon: GPO / OU structure enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1615</id>
    </mitre>
    <group>recon,discovery,ad_enumeration,</group>
  </rule>

  <!-- ============================================================
       SECTION 3 — Network Reconnaissance
       ============================================================ -->

  <!-- Network share enumeration -->
  <rule id="100420" level="10">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-SmbShare|Get-WmiObject Win32_Share|net share|net view</field>
    <options>no_full_log</options>
    <description>PS Recon: Network share enumeration on $(win.system.computer)</description>
    <mitre>
      <id>T1135</id>
    </mitre>
    <group>recon,discovery,network_recon,</group>
  </rule>

  <!-- Port scanning -->
  <rule id="100421" level="12">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Test-NetConnection|Invoke-Portscan|1..1024|foreach.*TcpClient</field>
    <options>no_full_log</options>
    <description>PS Recon: Port scanning activity detected on $(win.system.computer)</description>
    <mitre>
      <id>T1046</id>
    </mitre>
    <group>recon,discovery,network_recon,port_scan,</group>
  </rule>

  <!-- DNS reconnaissance -->
  <rule id="100422" level="8">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Resolve-DnsName|nslookup|\[System\.Net\.Dns\]</field>
    <options>no_full_log</options>
    <description>PS Recon: DNS reconnaissance on $(win.system.computer)</description>
    <mitre>
      <id>T1018</id>
    </mitre>
    <group>recon,discovery,network_recon,</group>
  </rule>

  <!-- ============================================================
       SECTION 4 — Credential Access
       ============================================================ -->

  <!-- Registry credential hunting -->
  <rule id="100430" level="13">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)HKLM.*password|HKCU.*password|reg query.*password|Get-ItemProperty.*password</field>
    <options>no_full_log</options>
    <description>PS Attack: Registry credential search on $(win.system.computer)</description>
    <mitre>
      <id>T1552.002</id>
    </mitre>
    <group>attack,credential_theft,</group>
  </rule>

  <!-- LSASS access attempts -->
  <rule id="100431" level="14">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)lsass|sekurlsa|Out-Minidump|MiniDumpWriteDump|procdump</field>
    <options>no_full_log</options>
    <description>PS Attack: LSASS memory access / credential dump attempt on $(win.system.computer)</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
    <group>attack,credential_theft,</group>
  </rule>

  <!-- SAM/NTDS credential dumping -->
  <rule id="100432" level="14">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)ntds.dit|sam.*system|reg save.*sam|reg save.*system|vssadmin</field>
    <options>no_full_log</options>
    <description>PS Attack: SAM/NTDS credential dump attempt on $(win.system.computer)</description>
    <mitre>
      <id>T1003.002</id>
    </mitre>
    <group>attack,credential_theft,</group>
  </rule>

  <!-- Credential files search -->
  <rule id="100433" level="12">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Get-ChildItem.*password|Select-String.*password|findstr.*password|Get-Content.*unattend</field>
    <options>no_full_log</options>
    <description>PS Attack: Credential file search on $(win.system.computer)</description>
    <mitre>
      <id>T1552.001</id>
    </mitre>
    <group>attack,credential_theft,</group>
  </rule>

  <!-- ============================================================
       SECTION 5 — Defense Evasion
       ============================================================ -->

  <!-- Disabling Windows Defender -->
  <rule id="100440" level="14">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Set-MpPreference|Disable-WindowsOptionalFeature|Add-MpPreference.*ExclusionPath|Set-MpPreference.*DisableRealtimeMonitoring</field>
    <options>no_full_log</options>
    <description>PS Attack: Windows Defender tampering on $(win.system.computer)</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
    <group>attack,defense_evasion,</group>
  </rule>

  <!-- Clearing event logs -->
  <rule id="100441" level="14">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Clear-EventLog|wevtutil.*cl|Remove-EventLog|Limit-EventLog</field>
    <options>no_full_log</options>
    <description>PS Attack: Event log clearing on $(win.system.computer)</description>
    <mitre>
      <id>T1070.001</id>
    </mitre>
    <group>attack,defense_evasion,log_cleared,</group>
  </rule>

  <!-- AMSI bypass attempts -->
  <rule id="100442" level="15">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)amsiInitFailed|AmsiScanBuffer|amsi.dll|Bypass-AMSI|amsiContext</field>
    <options>no_full_log</options>
    <description>PS Attack: AMSI bypass attempt on $(win.system.computer)</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
    <group>attack,defense_evasion,amsi_bypass,</group>
  </rule>

  <!-- Execution policy bypass -->
  <rule id="100443" level="10">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Set-ExecutionPolicy Bypass|Set-ExecutionPolicy Unrestricted|ExecutionPolicy Bypass</field>
    <options>no_full_log</options>
    <description>PS Attack: PowerShell execution policy bypass on $(win.system.computer)</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <group>attack,defense_evasion,</group>
  </rule>

  <!-- ============================================================
       SECTION 6 — Persistence
       ============================================================ -->

  <!-- Scheduled task creation -->
  <rule id="100450" level="12">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Register-ScheduledTask|New-ScheduledTask|schtasks.*create</field>
    <options>no_full_log</options>
    <description>PS Persistence: Scheduled task creation on $(win.system.computer)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>attack,persistence,</group>
  </rule>

  <!-- Registry run key persistence -->
  <rule id="100451" level="12">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)HKLM.*Run|HKCU.*Run|New-ItemProperty.*CurrentVersion.*Run|Set-ItemProperty.*CurrentVersion.*Run</field>
    <options>no_full_log</options>
    <description>PS Persistence: Registry run key modification on $(win.system.computer)</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
    <group>attack,persistence,</group>
  </rule>

  <!-- WMI event subscription persistence -->
  <rule id="100452" level="13">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)New-CimInstance.*__EventFilter|New-CimInstance.*__EventConsumer|New-CimInstance.*__FilterToConsumerBinding|Set-WmiInstance.*__EventFilter</field>
    <options>no_full_log</options>
    <description>PS Persistence: WMI event subscription created on $(win.system.computer)</description>
    <mitre>
      <id>T1546.003</id>
    </mitre>
    <group>attack,persistence,wmi_persistence,</group>
  </rule>

  <!-- ============================================================
       SECTION 7 — Lateral Movement
       ============================================================ -->

  <!-- Remote execution via WMI -->
  <rule id="100460" level="12">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Invoke-WmiMethod|Get-WmiObject.*-ComputerName|Invoke-CimMethod.*-ComputerName</field>
    <options>no_full_log</options>
    <description>PS Lateral Movement: Remote WMI execution from $(win.system.computer)</description>
    <mitre>
      <id>T1047</id>
    </mitre>
    <group>attack,lateral_movement,</group>
  </rule>

  <!-- Remote execution via PSRemoting -->
  <rule id="100461" level="12">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Invoke-Command.*-ComputerName|New-PSSession|Enter-PSSession|New-CimSession.*-ComputerName</field>
    <options>no_full_log</options>
    <description>PS Lateral Movement: PSRemoting / remote session from $(win.system.computer)</description>
    <mitre>
      <id>T1021.006</id>
    </mitre>
    <group>attack,lateral_movement,</group>
  </rule>

  <!-- Pass-the-hash / token manipulation -->
  <rule id="100462" level="14">
    <if_sid>91802</if_sid>
    <field name="win.system.eventID">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)Invoke-TokenManipulation|Invoke-RunAs|Add-ObjectAcl|Invoke-ACLScanner</field>
    <options>no_full_log</options>
    <description>PS Attack: Token manipulation / ACL abuse on $(win.system.computer)</description>
    <mitre>
      <id>T1134</id>
    </mitre>
    <group>attack,lateral_movement,privilege_escalation,</group>
  </rule>

</group>