# 102-registry_integrity_windows.xml
<group name="windows,registry_integrity,syscheck,">

  <!-- ============================================================
       SECTION 1 — Persistence (Run Keys)
  ============================================================ -->

  <!-- Run / RunOnce key modified — classic persistence -->
  <rule id="100500" level="12">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">CurrentVersion\\Run</field>
    <description>Registry: Run/RunOnce key changed — possible persistence on $(agent.name)</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
    <group>registry,persistence,</group>
  </rule>


  <!-- ============================================================
       SECTION 2 — Credential Protection (LSA / WDigest)
  ============================================================ -->

  <!-- LSA registry key modified — RunAsPPL, auth packages, notification DLLs -->
  <rule id="100501" level="13">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path" type="pcre2">Control\\Lsa($|\\)</field>
    <description>Registry: LSA key modified — credential protection settings changed on $(agent.name)</description>
    <mitre>
      <id>T1556.002</id>
    </mitre>
    <group>registry,credential_access,defense_evasion,</group>
  </rule>

  <!-- WDigest modified — plaintext credential caching may be re-enabled -->
  <rule id="100502" level="14">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">SecurityProviders\\WDigest</field>
    <description>Registry: WDigest key modified — plaintext credential caching may be enabled on $(agent.name)</description>
    <mitre>
      <id>T1556.002</id>
    </mitre>
    <group>registry,credential_access,</group>
  </rule>

  <!-- SSP/Security Provider list — injected DLLs load into LSASS -->
  <rule id="100503" level="14">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Control\\Lsa\\OSConfig</field>
    <description>Registry: LSA OSConfig modified — Security Support Provider list may have changed on $(agent.name)</description>
    <mitre>
      <id>T1547.005</id>
    </mitre>
    <group>registry,persistence,credential_access,</group>
  </rule>


  <!-- ============================================================
       SECTION 3 — Defense Evasion
  ============================================================ -->

  <!-- Windows Defender policy key modified -->
  <rule id="100510" level="14">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Policies\\Microsoft\\Windows Defender</field>
    <description>Registry: Windows Defender policy key modified on $(agent.name) — protection may be disabled</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
    <group>registry,defense_evasion,</group>
  </rule>

  <!-- Windows Defender real-time settings modified.
       Excludes the Policies subtree so this does not double-alert with 100510
       (which owns Policies\Microsoft\Windows Defender). -->
  <rule id="100511" level="14">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Microsoft\\Windows Defender</field>
    <field name="syscheck.path" negate="yes" type="pcre2">Policies\\Microsoft\\Windows Defender</field>
    <description>Registry: Windows Defender settings modified on $(agent.name)</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
    <group>registry,defense_evasion,</group>
  </rule>

  <!-- UAC policy modified — disabling prompts is a common escalation step -->
  <rule id="100512" level="13">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Policies\\System</field>
    <description>Registry: UAC/system policy key modified on $(agent.name)</description>
    <mitre>
      <id>T1548.002</id>
    </mitre>
    <group>registry,defense_evasion,privilege_escalation,</group>
  </rule>

  <!-- Audit policy registry key modified — possible attempt to blind logging -->
  <rule id="100513" level="14">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">SECURITY\\Policy\\PolAdtEv</field>
    <description>Registry: Audit policy registry key modified on $(agent.name) — logging may be reduced</description>
    <mitre>
      <id>T1562.002</id>
    </mitre>
    <group>registry,defense_evasion,</group>
  </rule>


  <!-- ============================================================
       SECTION 4 — Backdoors / Code Injection
  ============================================================ -->

  <!-- Image File Execution Options — debugger hijacking (accessibility backdoor) -->
  <rule id="100520" level="15">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Image File Execution Options</field>
    <description>Registry: IFEO key modified on $(agent.name) — possible debugger hijacking or accessibility backdoor</description>
    <mitre>
      <id>T1546.012</id>
    </mitre>
    <group>registry,persistence,privilege_escalation,</group>
  </rule>

  <!-- Winlogon modified — shell or credential provider replacement -->
  <rule id="100521" level="14">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">CurrentVersion\\Winlogon</field>
    <description>Registry: Winlogon key modified on $(agent.name) — shell or credential provider may have changed</description>
    <mitre>
      <id>T1547.004</id>
    </mitre>
    <group>registry,persistence,</group>
  </rule>

  <!-- AppInit_DLLs — DLL injected into every user-mode process -->
  <rule id="100522" level="15">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Windows NT\\CurrentVersion\\Windows</field>
    <description>Registry: AppInit_DLLs key modified on $(agent.name) — DLL injection into user-mode processes possible</description>
    <mitre>
      <id>T1546.010</id>
    </mitre>
    <group>registry,persistence,</group>
  </rule>


  <!-- ============================================================
       SECTION 5 — Services
  ============================================================ -->

  <!-- New service created or existing service modified -->
  <rule id="100530" level="11">
    <if_sid>550,553</if_sid>
    <field name="syscheck.path">CurrentControlSet\\Services</field>
    <description>Registry: Services key modified on $(agent.name) — service created or changed</description>
    <mitre>
      <id>T1543.003</id>
    </mitre>
    <group>registry,persistence,</group>
  </rule>

  <!-- Service deleted -->
  <rule id="100531" level="10">
    <if_sid>554</if_sid>
    <field name="syscheck.path">CurrentControlSet\\Services</field>
    <description>Registry: Services key deleted on $(agent.name) — service removed</description>
    <mitre>
      <id>T1543.003</id>
    </mitre>
    <group>registry,</group>
  </rule>


  <!-- ============================================================
       SECTION 6 — Remote Access / RDP
  ============================================================ -->

  <!-- RDP enable/disable or port change -->
  <rule id="100535" level="12">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Control\\Terminal Server</field>
    <description>Registry: Terminal Server/RDP key modified on $(agent.name) — RDP may have been enabled or reconfigured</description>
    <mitre>
      <id>T1021.001</id>
    </mitre>
    <group>registry,lateral_movement,</group>
  </rule>


  <!-- ============================================================
       SECTION 7 — PowerShell Execution Policy
  ============================================================ -->

  <!-- PowerShell execution policy lowered via registry -->
  <rule id="100540" level="12">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Microsoft\\PowerShell</field>
    <description>Registry: PowerShell execution policy key modified on $(agent.name)</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <group>registry,defense_evasion,</group>
  </rule>


  <!-- ============================================================
       SECTION 8 — DOMAIN CONTROLLER SPECIFIC
       Apply only when agent is in windows-dc group
  ============================================================ -->

  <!-- Print Spooler re-enabled on DC (should be disabled post-PrintNightmare) -->
  <rule id="100550" level="14">
    <if_sid>550,553</if_sid>
    <field name="syscheck.path">CurrentControlSet\\Services\\Spooler</field>
    <description>Registry: Print Spooler service key modified on DC $(agent.name) — verify spooler is not being re-enabled</description>
    <mitre>
      <id>T1574</id>
    </mitre>
    <group>registry,dc,defense_evasion,</group>
  </rule>

  <!-- Netlogon parameters changed — secure channel, DC locator behaviour -->
  <rule id="100551" level="12">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Services\\Netlogon\\Parameters</field>
    <description>Registry: Netlogon parameters modified on DC $(agent.name)</description>
    <mitre>
      <id>T1562</id>
    </mitre>
    <group>registry,dc,</group>
  </rule>

  <!-- NTDS parameters changed — AD database path or replication config -->
  <rule id="100552" level="13">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Services\\NTDS\\Parameters</field>
    <description>Registry: NTDS parameters modified on DC $(agent.name) — AD database config may have changed</description>
    <mitre>
      <id>T1562</id>
    </mitre>
    <group>registry,dc,</group>
  </rule>

  <!-- KDC / Kerberos settings changed — encryption types, ticket lifetime -->
  <rule id="100553" level="13">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Services\\Kdc</field>
    <description>Registry: KDC/Kerberos settings modified on DC $(agent.name) — ticket policy or encryption types may have changed</description>
    <mitre>
      <id>T1558</id>
    </mitre>
    <group>registry,dc,</group>
  </rule>

  <!-- W32Time config changed — Kerberos is time-sensitive, time skew = auth failure -->
  <rule id="100554" level="11">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Services\\W32Time</field>
    <options>no_full_log</options>
    <description>Registry: Windows Time service config modified on DC $(agent.name)</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
    <group>registry,dc,</group>
  </rule>

  <!-- DNS Server parameters changed -->
  <rule id="100555" level="11">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Services\\DNS\\Parameters</field>
    <options>no_full_log</options>
    <description>Registry: DNS Server parameters modified on DC $(agent.name)</description>
    <mitre>
      <id>T1584.002</id>
    </mitre>
    <group>registry,dc,</group>
  </rule>

  <!-- DFSR parameters changed — SYSVOL replication health -->
  <rule id="100556" level="11">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Services\\DFSR\\Parameters</field>
    <options>no_full_log</options>
    <description>Registry: DFSR parameters modified on DC $(agent.name) — SYSVOL replication config may have changed</description>
    <mitre>
      <id>T1490</id>
    </mitre>
    <group>registry,dc,</group>
  </rule>

  <!-- SMB server parameters — EternalBlue / PrintNightmare mitigations -->
  <rule id="100557" level="12">
    <if_sid>550,553,554</if_sid>
    <field name="syscheck.path">Services\\LanmanServer\\Parameters</field>
    <description>Registry: SMB server parameters modified on DC $(agent.name) — verify SMB hardening settings are intact</description>
    <mitre>
      <id>T1562</id>
    </mitre>
    <group>registry,dc,defense_evasion,</group>
  </rule>

</group>
